# APRS Services Loop/SERVICE-DAEMON/APRSHostStart.service
[Unit]
Description=Service to handle APRSHost

[Service]
ExecStart=/usr/bin/APRSHostStart.sh

[Install]
WantedBy=multi-user.target
